# other-configs/mirrors-and-hosts/rate-mirrors.timer
[Unit]
Description=rate-arch-mirrors
After=network.target

[Timer]
OnBootSec=1 m 
Unit=rate-mirrors.service

[Install]
WantedBy=multi-user.target
~
